# layout/Library/LaunchDaemons/com.qianmian.vcamreceiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qianmian.vcamreceiver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/var/jb/usr/local/bin/vcamreceiverd_launch.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/mobile/Library/VCam/vcamreceiver.log</string>
    <key>StandardErrorPath</key>
    <string>/var/mobile/Library/VCam/vcamreceiver.err</string>
</dict>
</plist>
